# ansible/roles/rq/templates/supervisor_rq.conf.j2
[program:rq]
command={{ rq_scripts_dir }}/rq_start

autostart=true
autorestart=true

user={{ rq_user }}

stdout_logfile={{ application_log_dir }}/rq.log
redirect_stderr = true