# ansible/roles/base/templates/supervisor_lfw.conf.j2
[program:logstash-forwarder]
command=/usr/local/bin/lfw_start

autostart=true
autorestart=true

user=root

stdout_logfile=/var/log/lfw.log
redirect_stderr = true
