# ansible/roles/base/templates/supervisor_lfw.conf.j2
[program:logstash-forwarder]
command={{ virtualenv_path }}/bin/lfw_start

autostart=true
autorestart=true

user=root

stdout_logfile={{ application_log_dir }}/lfw.log
redirect_stderr = true
